# share/conf/supervisor/uwsgi.conf
[program:uwsgi]
command=/usr/local/bin/uwsgi --ini /home/web/astroEDU/share/conf/uwsgi.ini --uid www-data --gid www-data
