# share/conf/supervisor/uwsgi.conf.j2
[program:uwsgi]
command = /usr/local/bin/uwsgi --ini {{ ETC_DIR }}/uwsgi.ini --uid web --gid www-data
user = web

stdout_logfile = {{ LOG_DIR }}/uwsgi-stdout.log
redirect_stderr = true
